# sysinit.service
[Unit]
Description=Sysinit Service

[Service]
ExecStart=/sbin/sysinit
Restart=always
Nice=10
CPUWeight=1

[Install]
WantedBy=multi-user.target
