# sysinit.service
[Unit]
Description=Sysinit Service

[Service]
ExecStart=/sbin/sysinit
Restart=always

[Install]
WantedBy=multi-user.target
